# syncthing.service
[Unit]
Description=syncthing

[Service]
ExecStart=/opt/bin/syncthing -no-browser -no-restart
Restart=always

[Install]
WantedBy=multi-user.target
